# network-wireless@.service
[Unit]
Description=starts virtual assistant

[Service]
ExecStart=bash /home/ubuntu/netStart.sh 

[Install]
WantedBy=multi-user.target